# VM/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan serve --host=0.0.0.0
#stdout_logfile=/var/log/php-fpm/stdout.log
#stdout_logfile_maxbytes=0
#stderr_logfile=/var/log/php-fpm/stderr.log
#stderr_logfile_maxbytes=0
autostart=true
autorestart=unexpected

[program:node]
directory=/var/www/html/client
autostart=true
command=npm run dev